# service/minecraftd.service
[Unit]
Description=Minecraft server %i

[Service]
Type=simple
User=minecraft
ExecStart=/usr/local/bin/mined %i

[Install]
WantedBy=multi-user.target
